# results/case3/LOW_output.json
[{"id": "daca7f93-5f2f-4ef8-8ff4-43f6e8228c02", "attack_steps": [{"id": "7d42ff22-93f9-4098-b716-115673b78c68", "name": "Command and Scripting Interpreter: Python as used by the malware", "description": "A Python script designated as \"Small Sieve,\" packaged utilizing PyInstaller, is employed by MuddyWater for command and scripting execution.  \n", "mitre_technique": {"id": "T1059.006", "name": "Command and scripting interpreter: python", "description": "Adversaries may abuse Python commands and scripts for execution. Python is a very popular scripting/programming language, with capabilities to perform many functions. Python can be executed interactively from the command-line (via the python.exe interpreter) or via scripts (.py) that can be written and distributed to different systems. Python code can also be compiled into binary executables."}, "pre-conditions": ["\nThe malware has access to the Python interpreter.\n", "\nThe malware must have the capability to execute code.\n", "\nA Python interpreter is available on the target system.\n"], "post-conditions": ["\nNew files (e.g., malicious scripts, payloads)\n", "\nBase64-encoded PowerShell commands\n", "\nMalware persistence\n", "\nCustom encoding scheme in network traffic\n", "\nSystem compromise\n", "\nNetwork traffic to adversary-controlled IP addresses\n", "\nHex-encoded data transfers\n", "\nCommand and control channel established\n", "\nData exfiltration\n", "\nAltered system logs\n", "\nModified registry keys\n"], "indicators": []}], "pre-conditions": ["\nThe malware has access to the Python interpreter.\n", "\nThe malware must have the capability to execute code.\n", "\nA Python interpreter is available on the target system.\n"], "post-conditions": ["\nNew files (e.g., malicious scripts, payloads)\n", "\nBase64-encoded PowerShell commands\n", "\nMalware persistence\n", "\nCustom encoding scheme in network traffic\n", "\nSystem compromise\n", "\nNetwork traffic to adversary-controlled IP addresses\n", "\nHex-encoded data transfers\n", "\nCommand and control channel established\n", "\nData exfiltration\n", "\nAltered system logs\n", "\nModified registry keys\n"]}, {"id": "f45c3713-6fe5-4858-87c0-bcaac5ad02d5", "attack_steps": [{"id": "1f6cbc34-9878-4fd1-b01f-9e75d5b03812", "name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder as used by the malware", "description": "Persistence is established by the malware through Boot or Logon Autostart Execution utilizing Registry Run Keys.  The malware inserts itself into the HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run registry key under the designation \"SystemTextEncoding\" (or potentially \"OutlookMicrosift\" in certain instances). Upon user logon, Windows parses this registry key and automatically executes the malicious code associated with it. This action ensures malware execution during each system boot or user login cycle. \n\n\n", "mitre_technique": {"id": "T1547.001", "name": "Boot or logon autostart execution: registry run keys / startup folder", "description": "Adversaries may achieve persistence by adding a program to a startup folder or referencing it with a Registry run key. Adding an entry to the \"run keys\" in the Registry or startup folder will cause the program referenced to be executed when a user logs in. These programs will be executed under the context of the user and will have the account's associated permissions level."}, "pre-conditions": ["\nThe system is running a supported operating system.\n", "\nKnowledge of the appropriate registry key paths for persistence.\n", "\nA user account exists with sufficient privileges to modify registry keys.\n", "\nA tool capable of modifying registry values (e.g., command line tools, scripting languages).\n", "\nAccess to the Windows Registry.\n"], "post-conditions": ["\nSuspicious files in AppData/Roaming directory\n", "\nBase64-encoded PowerShell commands\n", "\nMalware persistence\n", "\nNetwork connections to command and control servers\n", "\nLaZagne output files containing stolen credentials\n", "\nCompromised system credentials\n", "\nSystem instability\n", "\nData exfiltration\n", "\nEmail logs with phishing links clicked\n", "\nRegistry Run keys modifications\n", "\nLogs of executed commands and scripts\n", "\nModified system files\n"], "indicators": ["Registry value name: HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\OutlookMicrosift", "Path: %AppData%\\OutlookMicrosift\\index.exe"]}], "pre-conditions": ["\nThe system is running a supported operating system.\n", "\nKnowledge of the appropriate registry key paths for persistence.\n", "\nA user account exists with sufficient privileges to modify registry keys.\n", "\nA tool capable of modifying registry values (e.g., command line tools, scripting languages).\n", "\nAccess to the Windows Registry.\n"], "post-conditions": ["\nSuspicious files in AppData/Roaming directory\n", "\nBase64-encoded PowerShell commands\n", "\nMalware persistence\n", "\nNetwork connections to command and control servers\n", "\nLaZagne output files containing stolen credentials\n", "\nCompromised system credentials\n", "\nSystem instability\n", "\nData exfiltration\n", "\nEmail logs with phishing links clicked\n", "\nRegistry Run keys modifications\n", "\nLogs of executed commands and scripts\n", "\nModified system files\n"]}, {"id": "4be78f8e-f4c5-4249-a0f6-4a0c68228073", "attack_steps": [{"id": "baeae286-935a-4409-8987-e2a6c920cfc5", "name": "Obfuscated Files or Information as used by the malware", "description": "The malware, Small Sieve, employs a two-stage obfuscation process to conceal its files and information.  \n\nHex byte swapping is utilized to rearrange the order of bytes within program strings and Telegram credentials according to a proprietary algorithm. This manipulation renders the data seemingly random and unintelligible without knowledge of the corresponding decryption key. Subsequently, the modified data undergoes encoding via an obfuscated variant of the Base64 algorithm, introducing an additional layer of complexity that hinders recognition and deciphering by security tools.  \n\n\n", "mitre_technique": {"id": "T1027.013", "name": "Obfuscated files or information: encrypted/encoded file", "description": "Adversaries may encrypt or encode files to obfuscate strings, bytes, and other specific patterns to impede detection. Encrypting and/or encoding file content aims to conceal malicious artifacts within a file used in an intrusion. Many other techniques, such as Software Packing, Steganography, and Embedded Payloads, share this same broad objective. Encrypting and/or encoding files could lead to a lapse in detection of static signatures, only for this malicious content to be revealed (i.e., Deobfuscate/Decode Files or Information) at the time of execution/use."}, "pre-conditions": ["\nAccess to a custom hex byte swapping encoding scheme.\n", "\nAccess to an obfuscated Base64 function.\n", "\nProgram strings and Telegram credentials exist within the malware.\n", "\nThe malware Small Sieve is present in the system.\n"], "post-conditions": ["\nObfuscated program strings\n", "\nLogs containing suspicious activity\n", "\nEncrypted files with unusual extensions\n", "\nMalware persistence\n", "\nCompromised system\n", "\nAltered web server logs\n", "\nTelegram API communication logs\n", "\nModified system configuration files\n", "\nJSON files containing C2 commands or results\n", "\nBase64 encoded data\n", "\nNetwork traffic to adversary-controlled IP addresses\n", "\nSystem instability\n", "\nModified system registry entries\n", "\nData exfiltration\n", "\nHex byte swapped data\n"], "indicators": ["Small Sieve sample (Filename: gram_app.exe)", "Small Sieve sample (Filename: index.exe)", "Path: %LocalAppData%\\MicrosoftWindowsOutlookDataPlus.txt", "Registry value name: HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\OutlookMicrosift", "Path: %AppData%\\OutlookMicrosift\\index.exe"]}, {"id": "bef3e989-3501-4ee9-b4f6-1ee45263dfae", "name": "Execution Guardrails as used by the malware", "description": "The malware's functionality is contingent upon the fulfillment of a specific prerequisite known as an \"execution guardrail.\"  This technique necessitates the provision of a designated command-line argument, in this instance, the word \"Platypus,\" for the Small Sieve payload, employed by the MuddyWater group, to execute successfully. This parameter serves as a safeguard, inhibiting inadvertent or unauthorized malware activation. \n", "mitre_technique": {"id": "T1480.001", "name": "Execution guardrails: environmental keying", "description": "Adversaries may environmentally key payloads or other features of malware to evade defenses and constraint execution to a specific target environment. Environmental keying uses cryptography to constrain execution or actions based on adversary supplied environment specific conditions that are expected to be present on the target. Environmental keying is an implementation of Execution Guardrails that utilizes cryptographic techniques for deriving encryption/decryption keys from specific types of values in a given computing environment."}, "pre-conditions": ["\nThe word \"Platypus\" must be included in the command line arguments when executing Small Sieve.\n", "\nThe command line interface is accessible.\n", "\nThe malware Small Sieve is present in the system.\n", "\nA system running a compatible operating system.\n"], "post-conditions": ["\nNetwork traffic to command and control servers\n", "\nMalware persistence\n", "\nUnusual file access patterns\n", "\nLog entries indicating suspicious activity\n", "\nSystem instability\n", "\nNewly created files with malicious code\n", "\nAltered registry keys\n", "\nHidden processes running in memory\n", "\nObfuscated data transmissions\n", "\nData exfiltration\n", "\nModified system settings\n", "\nCompromised systems\n"], "indicators": ["Small Sieve sample (Filename: gram_app.exe)", "Small Sieve sample (Filename: index.exe)", "Path: %AppData%\\OutlookMicrosift\\index.exe"]}, {"id": "0784a0e2-6511-4f81-bd36-bf8152889a46", "name": "Masquerading: Match Legitimate Name or Location as used by the malware", "description": "Malicious filenames mimicking legitimate Microsoft products, such as \"Microsift\" and \"Outlook,\" are employed by the Small Sieve malware to facilitate concealment. This tactic aims to integrate seamlessly with existing files, thereby hindering detection through casual scrutiny.  \n", "mitre_technique": {"id": "T1036.005", "name": "Masquerading: match legitimate name or location", "description": "Adversaries may match or approximate the name or location of legitimate files or resources when naming/placing them. This is done for the sake of evading defenses and observation. This may be done by placing an executable in a commonly trusted directory (ex: under System32) or giving it the name of a legitimate, trusted program (ex: svchost.exe). In containerized environments, this may also be done by creating a resource in a namespace that matches the naming convention of a container pod or cluster. Alternatively, a file or container image name given may be a close approximation to legitimate programs/images or something innocuous."}, "pre-conditions": ["\nThe malware has knowledge of legitimate software names like Microsoft and Outlook.\n", "\nAbility to modify filenames.\n", "\nAccess to a file system.\n", "\nThe malware has access to a file system.\n"], "post-conditions": ["\nCompromised user credentials\n", "\nMalicious log files\n", "\nSuspicious network traffic (e.g., communication with command and control servers)\n", "\nTelegram API logs\n", "\nMalware persistence\n", "\nNetwork reconnaissance\n", "\nPresence of malware executables\n", "\nModified Outlook/Microsoft files with altered spellings (\"Microsift\")\n", "\nUnusual process activity\n", "\nModified or deleted files\n", "\nSystem instability\n", "\nAltered system registry entries\n", "\nCached credentials\n", "\nData exfiltration\n"], "indicators": ["Path: %LocalAppData%\\MicrosoftWindowsOutlookDataPlus.txt", "Registry value name: HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\OutlookMicrosift", "Path: %AppData%\\OutlookMicrosift\\index.exe"]}], "pre-conditions": ["\nAccess to a custom hex byte swapping encoding scheme.\n", "\nAccess to an obfuscated Base64 function.\n", "\nProgram strings and Telegram credentials exist within the malware.\n", "\nThe malware Small Sieve is present in the system.\n", "\nThe word \"Platypus\" must be included in the command line arguments when executing Small Sieve.\n", "\nThe command line interface is accessible.\n", "\nThe malware Small Sieve is present in the system.\n", "\nA system running a compatible operating system.\n", "\nThe malware has knowledge of legitimate software names like Microsoft and Outlook.\n", "\nAbility to modify filenames.\n", "\nAccess to a file system.\n", "\nThe malware has access to a file system.\n"], "post-conditions": ["\nObfuscated program strings\n", "\nLogs containing suspicious activity\n", "\nEncrypted files with unusual extensions\n", "\nMalware persistence\n", "\nCompromised system\n", "\nAltered web server logs\n", "\nTelegram API communication logs\n", "\nModified system configuration files\n", "\nJSON files containing C2 commands or results\n", "\nBase64 encoded data\n", "\nNetwork traffic to adversary-controlled IP addresses\n", "\nSystem instability\n", "\nModified system registry entries\n", "\nData exfiltration\n", "\nHex byte swapped data\n", "\nNetwork traffic to command and control servers\n", "\nMalware persistence\n", "\nUnusual file access patterns\n", "\nLog entries indicating suspicious activity\n", "\nSystem instability\n", "\nNewly created files with malicious code\n", "\nAltered registry keys\n", "\nHidden processes running in memory\n", "\nObfuscated data transmissions\n", "\nData exfiltration\n", "\nModified system settings\n", "\nCompromised systems\n", "\nCompromised user credentials\n", "\nMalicious log files\n", "\nSuspicious network traffic (e.g., communication with command and control servers)\n", "\nTelegram API logs\n", "\nMalware persistence\n", "\nNetwork reconnaissance\n", "\nPresence of malware executables\n", "\nModified Outlook/Microsoft files with altered spellings (\"Microsift\")\n", "\nUnusual process activity\n", "\nModified or deleted files\n", "\nSystem instability\n", "\nAltered system registry entries\n", "\nCached credentials\n", "\nData exfiltration\n"]}, {"id": "346d931d-6c3a-4f26-b86d-fdd00d036111", "attack_steps": [{"id": "3ac14fa0-de28-4eba-96b0-928e838fa9d1", "name": "Application Layer Protocol: Web Protocols as used by the malware", "description": "Application layer protocol operations, specifically web protocols, are executed by Small Sieve malware via utilization of the Telegram Bot Application Programming Interface (API) over HTTPS for both beaconing and tasking functionalities.  \n", "mitre_technique": {"id": "T1071.001", "name": "Application layer protocol: web protocols", "description": "Adversaries may communicate using application layer protocols associated with web traffic to avoid detection/network filtering by blending in with existing traffic. Commands to the remote system, and often the results of those commands, will be embedded within the protocol traffic between the client and server."}, "pre-conditions": ["\nA connection to the internet is available.\n", "\nThe malware has a valid Telegram token.\n", "\nThe Telegram Bot API is accessible.\n", "\nThe malware is running.\n", "\nThe malware has the capability to send and receive HTTP requests.\n"], "post-conditions": ["\nTelegram Bot API usage logs\n", "\nAltered registry settings\n", "\nMalware persistence\n", "\nDNS requests to malicious domains\n", "\nNew user accounts\n", "\nLog entries indicating suspicious activity\n", "\nSystem instability\n", "\nUnusual network traffic (e.g., communication with Telegram API)\n", "\nHidden files or folders containing malware code\n", "\nProxy server logs\n", "\nData exfiltration\n", "\nCompromised systems\n", "\nModified system files\n"], "indicators": []}, {"id": "cba42d89-baa6-4670-b1f6-6c25998ffde2", "name": "Data Encoding: Non-Standard Encoding as used by the malware", "description": "The tasking traffic transmitted between the malware and its command-and-control (C2) server undergoes obfuscation through a custom hex byte swapping encoding scheme.  \n", "mitre_technique": {"id": "T1132.002", "name": "Data encoding: non-standard encoding", "description": "Adversaries may encode data with a non-standard data encoding system to make the content of command and control traffic more difficult to detect. Command and control (C2) information can be encoded using a non-standard data encoding system that diverges from existing protocol specifications. Non-standard data encoding schemes may be based on or related to standard data encoding schemes, such as a modified Base64 encoding for the message body of an HTTP request."}, "pre-conditions": ["\nThe malware is running.\n", "\nThe Small Sieve tool is available within the malware's infrastructure.\n", "\nA communication channel with the C2 server exists.\n", "\nA system infected with MuddyWater malware.\n"], "post-conditions": ["\nNew log entries indicating suspicious activity.\n", "\nTelegram Bot API interaction logs.\n", "\nCompromised system with malware installed.\n", "\nEncrypted data files containing stolen information.\n", "\nAltered registry settings related to malware execution and persistence.\n", "\nNetwork traffic logs showing communication with adversary-controlled IP addresses.\n", "\nPresence of Small Sieve payload on the system.\n", "\nData exfiltration to adversary-controlled IP address.\n", "\nModified system files (e.g., DLLs, scripts).\n"], "indicators": []}], "pre-conditions": ["\nA connection to the internet is available.\n", "\nThe malware has a valid Telegram token.\n", "\nThe Telegram Bot API is accessible.\n", "\nThe malware is running.\n", "\nThe malware has the capability to send and receive HTTP requests.\n", "\nThe malware is running.\n", "\nThe Small Sieve tool is available within the malware's infrastructure.\n", "\nA communication channel with the C2 server exists.\n", "\nA system infected with MuddyWater malware.\n"], "post-conditions": ["\nTelegram Bot API usage logs\n", "\nAltered registry settings\n", "\nMalware persistence\n", "\nDNS requests to malicious domains\n", "\nNew user accounts\n", "\nLog entries indicating suspicious activity\n", "\nSystem instability\n", "\nUnusual network traffic (e.g., communication with Telegram API)\n", "\nHidden files or folders containing malware code\n", "\nProxy server logs\n", "\nData exfiltration\n", "\nCompromised systems\n", "\nModified system files\n", "\nNew log entries indicating suspicious activity.\n", "\nTelegram Bot API interaction logs.\n", "\nCompromised system with malware installed.\n", "\nEncrypted data files containing stolen information.\n", "\nAltered registry settings related to malware execution and persistence.\n", "\nNetwork traffic logs showing communication with adversary-controlled IP addresses.\n", "\nPresence of Small Sieve payload on the system.\n", "\nData exfiltration to adversary-controlled IP address.\n", "\nModified system files (e.g., DLLs, scripts).\n"]}]